# utils/layer1/centos-supervisord/container-files/etc/supervisor.d/munged.conf
[program:munged]
command=/usr/sbin/munged -F
user=munge
autostart=true
