# utils/layer1/centos-supervisord/container-files/etc/supervisor.d/rsyslogd.conf
[program:rsyslogd]
command=/usr/sbin/rsyslogd -n
autostart=true
