# utils/layer1/centos-supervisord/container-files/etc/supervisor.d/slurmd.conf
[program:slurmd]
command=/usr/sbin/slurmd -Dvvv
user=root
autostart=true
